# systemd/.config/systemd/user/redshift.service
[Unit]
Description=redshift

[Service]
Type=simple
ExecStart=/usr/bin/redshift

[Install]
WantedBy=default.target
