# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Removable disk automounter using udisks

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
