# systemd/.config/systemd/user/xbanish.service
[Unit]
Description=xbanish

[Service]
Type=simple
ExecStart=/usr/local/bin/xbanish

[Install]
WantedBy=default.target
